# .config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=10m
OnUnitActiveSec=45m
Unit=mbsync.service

[Install]
WantedBy=timers.target
